# systemd/lists-watcher.path
[Path]
PathModified=/home/droundy/thing-lists

[Install]
WantedBy=multi-user.target
