# systemd/lists-watcher.path
[Path]
PathModified=/home/droundy/lists/target/release/thing-lists

[Install]
WantedBy=multi-user.target
