# systemd/lists-watcher.path
[Path]
PathModified=/root/lists/.git

[Install]
WantedBy=multi-user.target
